# scripts/vault-watcher/vault-watcher.service
[Unit]
Description=AI Stack Vault File Watcher
After=network.target n8n-ai-stack.service

[Service]
Type=simple
User=root
WorkingDirectory=/mnt/user/appdata/ai_stack/scripts/vault-watcher
Environment="VAULT_DIR=/mnt/user/appdata/ai_stack/vault"
Environment="N8N_WEBHOOK=http://n8n-ai-stack:5678/webhook/reembed-file"
Environment="POSTGRES_HOST=postgres-ai-stack"
Environment="POSTGRES_PORT=5432"
Environment="POSTGRES_DB=aistack"
Environment="POSTGRES_USER=aistack_user"
Environment="POSTGRES_PASSWORD=changeme"
Environment="DEBOUNCE_SECONDS=5"
ExecStart=/mnt/user/appdata/ai_stack/scripts/vault-watcher/watch-vault.sh
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
